# supervise/redis.conf
[program:redis]
command=redis-server
priority=1
